# packaging/debian/dd-process-agent.prerm
#!/bin/sh

echo "Removing dd-process-agent from auto start on boot..."
update-rc.d -f dd-process-agent remove

echo "Trying to stop dd-process-agent with init.d"
/etc/init.d/dd-process-agent stop || true
exit 0
