# packaging/rpm/dd-process-agent.prerm
#!/bin/sh


echo "Removing dd-process-agent to start on boot..."
chkconfig --del dd-process-agent

echo "Trying to stop dd-process-agent with init.d"
/etc/init.d/dd-process-agent stop || true
exit 0
